# etc/cron.d/detect_passwd_change
# Detect any changes for file /etc/password every 15 minutes
*/15 * * * * root /usr/bin/detect_file_change /etc/passwd
